# master/files/etc/cron.d/retcon
MAILTO=bert@cyso.nl
1-59/15 * * * *	ubuntu bash /opt/scripts/update-retcon-stats 2>&1 | logger -p daemon.info -t retcon-runner
